# deploy/btc-bot.service
[Unit]
Description=BTC/MSTR ETF Signal Bot (worker)

[Service]
Type=simple
# %h = kullanici home dizini. Proje ~/MSTRadar altinda ise hicbir sey degistirmen GEREKMEZ.
# Baska bir klasore koyduysan iki yoldaki "MSTRadar" kismini guncelle.
WorkingDirectory=%h/MSTRadar
ExecStart=%h/MSTRadar/venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
